# kestrel/kestrel.upstart
description     "Kestrel"

start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn
respawn limit 10 5

console log
chdir /usr/share/kestrel-%%VERSION%%

script
  if [ -f /etc/default/kestrel ]; then
    . /etc/default/kestrel
  fi

  ulimit -n $KESTREL_FD_LIMIT
  ulimit -c unlimited

  /usr/bin/java $KESTREL_JAVA_OPTS -jar kestrel-%%VERSION%%.jar
end script
